# assistant-bot.service
[Unit]
Description=Personal Assistant Telegram Bot
After=network.target mysql.service
Wants=mysql.service

[Service]
User=www-data
WorkingDirectory=/opt/assistant_bot
ExecStart=/opt/assistant_bot/venv/bin/python main.py
Restart=always
RestartSec=10
EnvironmentFile=/opt/assistant_bot/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=assistant-bot

[Install]
WantedBy=multi-user.target
